# net-vpn/onioncat/files/garlicat.initd
#!/sbin/openrc-run
# Copyright 1999-2017 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

: ${UG:="onioncat:onioncat"}
: ${PIDFILE:="/run/gcat.pid"}

name="GaliCat"
command="/usr/bin/gcat"
command_args="-P ${PIDFILE} ${GCAT_OPTS}"
command_background=yes
pidfile="${PIDFILE}"
start_stop_daemon_args="-u ${UG}"

depend() {
	use net
	need i2p i2pd
}
