# x11-apps/radeon-profile-daemon/files/radeon-profile-daemon.initd
#!/sbin/openrc-run
# Copyright 1999-2017 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

pidfile="/run/radeon-profile-daemon.pid"
command="/usr/sbin/radeon-profile-daemon"
command_background=yes

depend() {
    need localmount
    after bootmisc
}
