# dev-db/influxdb/files/influxdb.initd-r1
#!/sbin/openrc-run
# Copyright 1999-2017 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

: ${INFLUXDB_USER:=influxdb}
: ${INFLUXDB_GROUP:=influxdb}
: ${INFLUXDB_CONFIGFILE:="/etc/influxdb/influxdb.conf"}
: ${INFLUXDB_LOGFILE:="/var/log/influxdb/influxd.log"}
: ${SSD_OPTS:="-1 ${INFLUXDB_LOGFILE} -2 ${INFLUXDB_LOGFILE}"}

name="InfluxDB daemon"
command="/usr/bin/influxd"
command_args="-config ${CONFIG}"
command_background=yes
pidfile="/run/${SVCNAME}.pid"
required_files="$INFLUXDB_CONFIGFILE"
start_stop_daemon_args="-u ${INFLUXDB_USER}
	-g ${INFLUXDB_GROUP} ${SSD_OPTS}"

depend() {
	need net
	after bootmisc
}
